# systemd/gh-zombie-reaper.service
[Unit]
Description=Kill stale gh CLI processes stuck in rate-limit retry loops

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gh-zombie-reaper.sh
