# systemd/hive-discord.service
[Unit]
Description=Hive Discord Bot
After=network-online.target hive.service
Wants=network-online.target
Requires=hive.service

[Service]
Type=simple
User=dev
Group=dev
WorkingDirectory=/tmp/hive/discord
EnvironmentFile=/etc/hive/discord.env
EnvironmentFile=-/etc/hive/supervisor.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node bot.js
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300
PrivateTmp=yes
ProtectSystem=strict
NoNewPrivileges=yes
ProtectHome=read-only
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hive-discord

[Install]
WantedBy=multi-user.target
