# systemd/hive-deploy.service
[Unit]
Description=Pull latest hive repo and sync scripts to /usr/local/bin
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=dev
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/local/bin/hive-deploy.sh
PrivateTmp=yes
ProtectSystem=strict
NoNewPrivileges=yes
ProtectHome=read-only
ReadWritePaths=/tmp/hive /usr/local/bin
StandardOutput=journal
StandardError=journal
