# systemd/hive-snapshot.timer
[Unit]
Description=Publish hive dashboard snapshot every 15 minutes

[Timer]
OnCalendar=*:00/15
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
