# systemd/hive@.service
[Unit]
Description=Hive agent %i (24/7 session with auto-restart)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dev
Group=dev
# Per-instance env file; default to /etc/hive/<instance>.env.
# %i is the instance name from hive@<name>.service.
EnvironmentFile=/etc/hive/%i.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/local/bin/supervisor.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
