# systemd/hive.service
[Unit]
Description=Hive agent (24/7 session with auto-restart)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dev
Group=dev
EnvironmentFile=/etc/hive/supervisor.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/local/bin/supervisor.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
