# systemd/hive-healthcheck.timer
[Unit]
Description=Run the hive healthcheck every 30 minutes

[Timer]
# First check 5 minutes after boot so the agent has time to come up.
OnBootSec=5min
# Then every 30 minutes (reduced from 20min to lower GH API rate limit pressure).
OnUnitActiveSec=30min
Unit=hive-healthcheck.service

[Install]
WantedBy=timers.target
