# systemd/kick-reviewer.service
[Unit]
Description=Kick reviewer agent to run a review pass
After=network-online.target

[Service]
Type=oneshot
User=dev
Environment=TMUX_BIN=tmux
ExecStart=/usr/local/bin/kick-agents.sh reviewer
StandardOutput=journal
StandardError=journal
