# systemd/kick-architect.service
[Unit]
Description=Kick architect agent to run an architecture/ideation pass
After=network-online.target

[Service]
Type=oneshot
User=dev
Environment=TMUX_BIN=tmux
ExecStart=/usr/local/bin/kick-agents.sh architect
StandardOutput=journal
StandardError=journal
