# systemd/kick-scanner.service
[Unit]
Description=Kick issue-scanner agent to run a scan pass
After=network-online.target

[Service]
Type=oneshot
User=dev
Environment=TMUX_BIN=tmux
ExecStart=/usr/local/bin/kick-agents.sh scanner
StandardOutput=journal
StandardError=journal
