# launchd/com.supervised-agent.scanner.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  supervised-agent scanner — macOS LaunchAgent
  Fires a standalone scanner script on a fixed schedule.

  This is the "hybrid" pattern: a deterministic bash scanner runs on a timer
  (no AI needed for scanning), writes state to SQLite, and triggers the AI
  agent only when there's actionable work.

  See examples/worker.sh.example for the scanner script pattern.
  See examples/kubestellar-fixer.md for a full case study.

  Install:
    cp launchd/com.supervised-agent.scanner.plist.example \
       ~/Library/LaunchAgents/com.supervised-agent.scanner.plist
    # Edit paths below
    launchctl load -w ~/Library/LaunchAgents/com.supervised-agent.scanner.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.supervised-agent.scanner.plist
    rm ~/Library/LaunchAgents/com.supervised-agent.scanner.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.supervised-agent.scanner</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/bash</string>
        <string>/Users/YOURUSER/.my-project-fix-loop/worker.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOURUSER</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- Fire every 15 minutes -->
    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/Users/YOURUSER/.my-project-fix-loop/scanner.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOURUSER/.my-project-fix-loop/scanner.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/YOURUSER</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
